# templates/rc.local.j2
#!/bin/bash

ip link set {{ listener_int }} promisc on
ip link set {{ listener_int }} up


exit 0